# deploy/puretrans-backend.service.example
[Unit]
Description=PureTrans Backend Service
After=network.target

[Service]
Type=simple
User=puretrans
Group=puretrans
WorkingDirectory=/opt/puretrans/backend
EnvironmentFile=/etc/puretrans/backend.env
ExecStart=/usr/bin/env uv run uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
